# notify-done.service
[Unit]
Description=Notify Done - Task completion notification daemon
Documentation=https://github.com/your-username/notify-done

[Service]
Type=simple
ExecStart=%h/.cargo/bin/nd daemon start
Environment=ND_DAEMON_CHILD=1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
